# root/var/service/faxq/run
#!/bin/sh

. /etc/sysconfig/i18n
export LANG
exec /usr/sbin/faxq -D 2>&1

